# tools/systemd/performance.service
[Unit]
Description=cpupower

[Service]
ExecStart=/usr/bin/cpupower frequency-set -g performance
Type=oneshot

[Install]
WantedBy=default.target
